# infra/systemd/ai-memory-sync.timer
[Unit]
Description=Run AI memory sync every 10 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
AccuracySec=30s
Persistent=true
Unit=ai-memory-sync.service

[Install]
WantedBy=timers.target
